# supervisord.conf
[supervisord]
nodaemon=true

[program:minidlna]
command=/usr/sbin/minidlnad -d
autostart=true
autorestart=true
